# infra/systemd/p2pool-backup.service
[Unit]
Description=P2Pool Dashboard database backup
Documentation=https://github.com/acquiredl/xmr-p2pool-dashboard
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/p2pool-dashboard

# Load env for BACKUP_REMOTE_URL, AWS credentials, etc.
EnvironmentFile=-/opt/p2pool-dashboard/.env

ExecStart=/bin/bash /opt/p2pool-dashboard/infra/scripts/pool-backup.sh

TimeoutStartSec=300
